# deploy/cron/crontab
* * * * * cd /app/backend && /usr/local/bin/python manage.py expire_stale_purchases >> /proc/1/fd/1 2>&1
